# scripts/com.ezhalah.wasalt-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ezhalah.wasalt-cleanup</string>
  <key>ProgramArguments</key>
  <array><string>/Users/yusufalnashwan/Downloads/design_handoff_ezhalah/ezhalah-app/scripts/wasalt-cleanup.sh</string></array>
  <key>StartInterval</key><integer>604800</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/tmp/com.ezhalah.wasalt-cleanup.out</string>
  <key>StandardErrorPath</key><string>/tmp/com.ezhalah.wasalt-cleanup.err</string>
  <key>Nice</key><integer>5</integer>
</dict>
</plist>
