# scripts/com.ezhalah.wasalt-sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- macOS launchd job: runs the Wasalt 4-hour sweep on this Mac.
       wasalt.sa geo-blocks GitHub's cloud IPs, so this is the working freshness path.
       Mac asleep → cycle skips (catches up on next wake). -->
  <key>Label</key>
  <string>com.ezhalah.wasalt-sweep</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/yusufalnashwan/Downloads/design_handoff_ezhalah/ezhalah-app/scripts/wasalt-sweep.sh</string>
  </array>

  <!-- Every 4 hours (14,400 sec). launchd fires the first run a few seconds after load. -->
  <key>StartInterval</key>
  <integer>14400</integer>

  <!-- Don't pile up missed cycles if the Mac was asleep — one catch-up is plenty. -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Logs (also tee'd into /tmp/wasalt-sweep.log by the script itself). -->
  <key>StandardOutPath</key>
  <string>/tmp/wasalt-sweep.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/wasalt-sweep.err</string>

  <!-- Nice to other processes; sweep is background hygiene. -->
  <key>Nice</key>
  <integer>5</integer>
</dict>
</plist>
